# monitor.service.default
[Unit]
Description=Controller for other services

[Service]
ExecStart=monitrage /etc/monitor.cfg
IgnoreSIGPIPE=false

[Install]
WantedBy=multi-user.target
Alias=monitor.service
