# monitor.service.default
[Unit]
Description=Controller for other services

[Service]
ExecStart=/usr/local/bin/monitrage /etc/monitor.cfg
WorkingDirectory=/tmp
IgnoreSIGPIPE=false

[Install]
WantedBy=multi-user.target
Alias=monitor.service
